# faq.service
[Unit]
Description=FAQ Backend Service
After=syslog.target network.target postgresql.service

[Service]
User=root
ExecStart=/usr/bin/java -Xmx512M -jar /root/faq-0.0.1-SNAPSHOT.jar
SuccessExitStatus=143
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target